# systemd/pcsc-over-http.socket
[Socket]
ListenStream=44602

[Install]
WantedBy=sockets.target