# systemd/pcsc-over-http.service
[Unit]
Description=PC/SC over HTTP service

[Service]
Type=simple
WorkingDirectory=/home/pi/.pcsc-over-http
ExecStart=/home/pi/.nvm/versions/node/v6.10.2/bin/node dist/index.js
StandardOutput=syslog+console
StandardError=syslog+console
SyslogIdentifier=pcsc-over-http
Environment=NODE_VERSION=6.10.2
Environment=PATH=/home/pi/.nvm/versions/node/v6.10.2/bin:/usr/bin:/usr/local/bin:/bin
Restart=on-abort
RestartSec=10

[Install]
WantedBy=multi-user.target