# systemd/pcsc-over-http.service
[Unit]
Description=PC/SC over HTTP service

[Service]
Type=simple
WorkingDirectory=/home/pi/.pcsc-over-http
ExecStart=/home/pi/.nvm/nvm-exec yarn serve
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=pcsc-over-http
User=nobody
Group=nobody
Environment=NODE_VERSION=6.10.2

[Install]
WantedBy=multi-user.target